# deploy/production/upstart.conf
start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 3 5

setuid dsnac
setgid dsnac

exec /home/dsnac/production/deploy/production/service_demon.sh